# image/meta-gcp-iot/recipes-gcp/mender-google-activation-agent/files/iot-activation.service
[Unit]
Description=waits for config from IoT core, then activates mender agent
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python /var/gcp/activate_agent.py

[Install]
WantedBy=multi-user.target
